# systemd/home-netops-tencent-firewall.timer
[Unit]
Description=Run home-netops Tencent firewall update periodically

[Timer]
OnBootSec=30s
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=home-netops-tencent-firewall.service

[Install]
WantedBy=timers.target
