# systemd/home-netops-tencent-firewall.service
[Unit]
Description=home-netops Tencent firewall update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-/etc/home-netops/home-netops.conf
ExecStart=/usr/local/lib/home-netops/firewall/tencent.sh
